# service/rsyncd/run
#!/bin/sh

if [ -e /var/run/rsyncd.pid ]; then
    rm /var/run/rsyncd.pid
fi

exec /usr/bin/rsync --no-detach --daemon --config /etc/rsyncd.conf
